# ha_automation_api/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

cd /usr/src/app
exec node src/server.js
